# master/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
/usr/sbin/php-fpm81 --nodaemonize --force-stderr
