# rootfs/etc/s6-overlay/s6-rc.d/init/up
#!/command/execlineb -P
with-contenv
/bin/init-container
